# resources/qemu/toils.service
[Unit]
Description=Toils Setup

[Service]
Type=simple
ExecStart=/usr/local/sbin/toils-everyboot.sh

[Install]
WantedBy=multi-user.target
