# resources/qemu/toils-setup.service
[Unit]
Description=Toils Setup

[Service]
Type=simple
ExecStart=/usr/local/sbin/toils-setup.sh

[Install]
WantedBy=multi-user.target
